# target/linux/imx8/base-files/etc/init.d/fiveg
#!/bin/sh /etc/rc.common

START=82

echo "fiveg..."
boot() {
	echo "fiveg boot..."
	/usr/bin/fiveg.sh &
}

start() {
	echo "fiveg start..."
}
